# templates/default/flyingduck.service.erb
[Unit]
Description = Flying Duck service (Arrow Flight Server with DuckDB).
After = syslog.target network.target docker.service <%= @local_dependencies %>
Requires = docker.service <%= @local_dependencies %>

[Service]
ExecStart = /usr/bin/docker run --rm --name %n \
      --memory=<%= node['flyingduck']['memory_gbs'] %>G \
      --hostname "0.0.0.0" \
      --network "host" \
      --volume "<%= @crypto_dir %>:/srv/hops/super_crypto/flyingduck" \
      --volume "<%= node['hops']['conf_dir_parent'] %>:/srv/hops/hadoop/etc" \
      --volume "<%= node['flyingduck']['logs'] %>:/usr/src/app/logs" \
      --volume "<%= node['flyingduck']['etc'] %>:/usr/src/app/etc" \
      --env "HADOOP_USER_NAME=<%= node['flyingduck']['user'] %>" \
      --env "BASE_DIR=/srv/hops" \
      --env "LIBHDFS_OPTS=-Dlog4j.configurationFile=/srv/hops/hadoop/etc/hadoop/log4j2.properties -Dhadoop.log.dir=/tmp -Dhadoop.root.logger=ERROR,console" \
      <%= @registry_image %> python src/server.py \
            --memory_limit=<%= node['flyingduck']['memory_gbs'] %> \
            --request_queue_timeout=<%= node['flyingduck']['request_timeout_sec'] %> \
            --query_result_batch_size=<%= node['flyingduck']['batch_size'] %> \
            --fqdn=<%= @flyingduck_fqdn %> \
            --port=<%= node['flyingduck']['port'] %> \
            --monitoring_port=<%= node['flyingduck']['monitoring'] %>

ExecStop = /usr/bin/docker stop flyingduck.service

Restart = on-failure
RestartSec = 2

[Install]
WantedBy = multi-user.target
